# crontab/crontab.example
* * * * * cd ~/coin && bash start.bash > /dev/null 2>/dev/null
